# meta-gewgaw/recipes-apps/gewgaw-collector/files/gewgaw-submit.timer
[Unit]
Description=Run gewgaw submit hourly

[Timer]
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
